# sound-rain-monitoring-py.service
# File Generated by service_maker.bash
[Unit]
Description=Script captures and outputs data from ADC and RG9 UART
After=multi-user.target
Conflicts=

[Service]
Type=simple
	
ExecStart=/usr/bin/python /home/pi/upstreamrain/sound-rain-monitoring.py
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
